# addons/gluon-ffp-xmlcollect/files/ffp-xmlcollect.cron
*/10 * * * * /usr/bin/ffp-xmlcollect.sh upload
*    * * * * /usr/bin/ffp-xmlcollect.sh collect
